# crontab
*/2 * * * * python3 /sheets2projectionlab.py >> /logfile.log 2>&1
