# crontab
0 16 * * * python /sheets2projectionlab.py >> /logfile